# 脚本/ck_mem_monitor.ini
[program:ck_mem_monitor]
process_name=%(program_name)s

command=/usr/bin/python3 -u /root/ck_mem_monitor.py

directory=/root

autostart=true
autorestart=true
startsecs=5
startretries=3

user=root

redirect_stderr=true
stdout_logfile=/var/log/ck_mem_monitor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
